# apps/DevOpsConsole/deploy/dev-coordinator.service
[Unit]
Description=DevCoordinator loopback control API
Documentation=file://%h/../DevCoordinator/skills/codex-dev-coordinator/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=holyglory
Group=holyglory
WorkingDirectory=/home/DevCoordinator
Environment=CODEX_AGENT_COORDINATOR_HOME=%h/.codex/agent-coordinator
ExecStart=/usr/bin/python3 %h/../DevCoordinator/skills/codex-dev-coordinator/scripts/dev_coordinator.py api serve --host 127.0.0.1 --port 29876 --token-file %h/.codex/agent-coordinator/api-token
Restart=on-failure
RestartSec=3
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=%h/.codex/agent-coordinator

[Install]
WantedBy=multi-user.target
