# apps/DevOpsConsole/deploy/devops-console.service
[Unit]
Description=DevOps Console - TLS edge, subdomain proxy, and control panel for vr.ae
Documentation=file://%h/../DevCoordinator/apps/DevOpsConsole/README.md
Requires=dev-coordinator.service
After=network-online.target dev-coordinator.service
Wants=network-online.target

[Service]
Type=simple
User=holyglory
Group=holyglory
WorkingDirectory=%h/../DevCoordinator/apps/DevOpsConsole
EnvironmentFile=%h/.config/devops-console/console.env
Environment=DEVCOORDINATOR_ROOT=/home/DevCoordinator
Environment=COORDINATOR_AUTOSTART=0
Environment=COORDINATOR_URL=http://127.0.0.1:29876
Environment=COORDINATOR_TOKEN_FILE=%h/.codex/agent-coordinator/api-token
Environment=CODEX_AGENT_COORDINATOR_HOME=%h/.codex/agent-coordinator
ExecStart=/usr/bin/node bin/devops-console.mjs
# SIGHUP reloads the TLS certificate from disk (after certbot renewal).
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3

# Bind 80/443 without root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

# The console never starts the coordinator. Runtime state and credentials are
# external to the Git checkout; TLS material remains read-only.
ProtectSystem=full
PrivateTmp=true
ReadWritePaths=%h/.local/state/devops-console
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
